# cookbooks/wosc-tabu/templates/default/supervisor.conf
[program:tabu]
command = /srv/tabu/deployment/bin/tabu-serve 7080

user = tabu

autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
