# cookbooks/wosc-prometheus/templates/default/supervisor/pushgateway.conf
[program:prometheus-push]
command = /srv/prometheus/bin/pushgateway

user = prometheus

redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
